# hosts-gen.service
[Unit]
Description=Generate /etc/hosts from several separate config files

[Service]
ExecStart=/usr/bin/hosts-gen
